# udev_rules/61-onacpower.rules
# Rule for when switching to AC Power
SUBSYSTEM=="power_supply", \
ATTR{online}=="0", \
RUN+="/usr/share/onBatteryLevelChanged.sh"
